# example/example.conf
[program:example]
command=/handler/start_actor.py ExampleActor
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true

