# serfnode/programs/docker_run.conf
[program:docker_run]
command={{DOCKER}} -H {{DOCKER_SOCKET}} run {{DOCKER_RUN}}
autostart=true
