# serfnode/programs/launcher.conf
[program:{{target}}]
command=/handler/launcher.py {{POS}} {{NAME}} {{ARGS}}
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stopsignal=INT
